# macos/launchd/com.pai.kernel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  PAI kernel LaunchAgent. Install via:
    pai-install-launchd install
  Uninstall:
    pai-install-launchd uninstall
  Status:
    pai-install-launchd status

  YOUR_HOME placeholders are substituted at install time.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pai.kernel</string>

    <key>ProgramArguments</key>
    <array>
        <string>YOUR_HOME/.pai/sbin/init</string>
    </array>

    <key>WorkingDirectory</key>
    <string>YOUR_HOME/.pai</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PAI_ROOT</key>
        <string>YOUR_HOME/.pai</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <!--
      Throttle restart attempts. Without this, a kernel that crashes during
      bootstrap will be re-spawned in a tight loop.
    -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>YOUR_HOME/.pai/var/log/kernel.out.log</string>

    <key>StandardErrorPath</key>
    <string>YOUR_HOME/.pai/var/log/kernel.err.log</string>

    <!--
      ProcessType=Interactive: scheduled like a foreground app, not throttled
      as a background daemon. The kernel does interactive work on the owner's
      behalf and should not be aggressively suspended by macOS.
    -->
    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
